# dtr5.socket
[Unit]
Description=django dtr5 socket

[Socket]
ListenStream=/run/dtr5.socket

[Install]
WantedBy=dtr5.target
